# scheduled-reboot-0.11/DEBIAN/prerm
#!/bin/bash
systemctl disable post-reboot.service
