# scheduled-reboot-0.15/DEBIAN/postinst
#!/bin/bash
systemctl enable post-reboot.service
